# wealthhub.service
[Unit]
Description=WealthHub Portfolio Management App
After=docker.service umbrel.service
Requires=docker.service
PartOf=umbrel.service

[Service]
Type=simple
User=umbrel
WorkingDirectory=/home/umbrel/umbrel/services/wealthhub
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=10

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/home/umbrel/umbrel/services/wealthhub/.env

[Install]
WantedBy=multi-user.target
